# deploy/thoth-hindsight.service
# thoth-hindsight.service — the standalone hindsight-api HTTP server.
#
# Hindsight is a REBUILDABLE derived index over the canonical vault (SPEC sections 8/15),
# never the store of record. thoth talks to it as a plain HTTP client (see
# src/thoth/hindsight.py): retain/recall/forget/reset are requests to
# http://127.0.0.1:8888, with the bank id as a path segment.
#
# This unit runs `hindsight-api` as a normal long-running foreground server, so systemd
# supervises it directly (Type=simple): a crash is visible, and `systemctl restart` works.
# It binds loopback only (--host 127.0.0.1); nothing external should reach it.
#
# Runs as the unprivileged `pkm` service user — like every thoth unit, it never touches
# root. The server's LLM/embedding config and DATABASE_URL are HINDSIGHT_API_* vars loaded
# from the shared EnvironmentFile (the same .env the other units use); no secret is inlined
# here. `<owner>` finalises the exact install paths during deployment.

[Unit]
Description=thoth hindsight-api semantic-index server (127.0.0.1:8888)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pkm
Group=pkm
# HOME and a pkm-readable CWD are load-bearing (the server stats ./.env via dotenv, so a
# root-owned CWD raises PermissionError). PATH must expose the hindsight-api entry point.
Environment=HOME=/home/pkm
Environment=PATH=/home/pkm/.local/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/home/pkm
# HINDSIGHT_API_* (LLM provider/model/key, DATABASE_URL, WORKER_ID) come from the shared
# .env — see deploy/.env.example. No secret is inlined in this unit.
EnvironmentFile=/home/pkm/.thoth/.env
ExecStart=/home/pkm/.local/bin/hindsight-api --host 127.0.0.1 --port 8888
Restart=on-failure
RestartSec=5
# Cold start can include Postgres connection + model warm-up.
TimeoutStartSec=180
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
