# deploy/thoth-slack.service
# thoth-slack.service — the unattended Slack Socket-Mode capture/retrieve daemon.
#
# This is the single long-running unit of the appliance (SPEC sections 2, 4, 13 and
# Appendix -> Backup/recovery step 6). It runs `thoth slack`, which blocks serving Slack
# Socket Mode. The cron-driven jobs (reindex, summary, config-backup) are NOT services —
# they are one-shot `thoth ...` invocations from the system crontab (see deploy/crontab).
#
# Unit name: the SPEC body (sections 2/4) calls this `pkm-slack.service` while the
# recovery checklist (Appendix step 6) names `thoth-slack.service`. The discrepancy is
# resolved here in favour of **thoth-slack.service** — it matches the app/repo name
# (`thoth`) and the operational recovery step that re-enables it. Install it under that
# name (`/etc/systemd/system/thoth-slack.service`).
#
# Secrets live only in ~/.thoth/.env (chmod 600, gitignored, never in any repo); the unit
# sources them via EnvironmentFile so no token is ever written into this tracked file. The
# two non-secret resolution vars are set inline per the SPEC (PKM_VAULT and the Obsidian
# vault name the app resolves before any file op).

[Unit]
Description=thoth PKM Slack daemon (capture/retrieve, Socket Mode)
After=network-online.target
Wants=network-online.target
# Hindsight runs as its own always-on unit (thoth-hindsight.service); this daemon only
# TALKS to it over 127.0.0.1:9277 and must never spawn it (ProtectHome=read-only below).
# Order after it and pull it in so a boot brings the index up first.
After=thoth-hindsight.service
Wants=thoth-hindsight.service

# Bound restart flapping (SPEC section 10 supervision, issue #15). With
# Restart=on-failure (below), systemd restarts the daemon on a crash; StartLimitBurst
# restarts within StartLimitIntervalSec then enters a permanent `failed` state instead of
# looping forever. The StartLimit* directives live in [Unit] (systemd moved them here from
# [Service]). How flapping SURFACES to the human: this box has no on-host monitoring, so
# the backstops are (1) the errors-to-Slack handler — the daemon posts each unhandled
# exception to SLACK_ALERT_CHANNEL before it exits (thoth.slack_app.serve_with_alerting),
# so a crash loop yields a burst of Slack alerts; and (2) the daily heartbeat — once the
# unit is `failed` it stops recording the `capture`/`push` liveness markers, so the next
# 07:00 digest's "still alive — last ingest/push at T" line goes stale, making the silence
# itself diagnostic. For a hard backstop, add an OnFailure= unit that pings a dead-man's
# switch (e.g. healthchecks.io) on entering `failed`; left out of the shipped unit to keep
# zero external deps, but documented here.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
# Dedicated unprivileged account for least privilege (SPEC section 12). Create the `pkm`
# user + group before enabling, or edit these two lines to an existing unprivileged user.
User=pkm
Group=pkm

# Non-secret resolution config (SPEC Appendix -> Vault identity). The app resolves the
# concrete absolute vault path from PKM_VAULT before any file op.
Environment=PKM_VAULT=/opt/pkm-vault
Environment=OBSIDIAN_VAULT_NAME=pkm-vault
Environment=THOTH_HOME=/home/pkm/.thoth

# Secrets (ANTHROPIC_API_KEY, SLACK_BOT_TOKEN, SLACK_APP_TOKEN, SLACK_ALLOWED_USERS,
# SLACK_SUMMARY_CHANNEL, EXA_API_KEY, FIRECRAWL_API_KEY, GEMINI_API_KEY, ...) come from the
# chmod-600 .env — never from this tracked unit file. Optional: set SLACK_ALERT_CHANNEL in
# .env to route errors-to-Slack alerts to a dedicated channel; unset, alerts fall back to the
# first SLACK_ALLOWED_USERS DM.
EnvironmentFile=/home/pkm/.thoth/.env

WorkingDirectory=/opt/thoth
ExecStart=/opt/thoth/.venv/bin/thoth slack
# Restart on a crash, paced by RestartSec, and bounded by the [Unit] StartLimit* so a
# wedged daemon enters `failed` (and stops the heartbeat) rather than spinning (issue #15).
Restart=on-failure
RestartSec=5

# Hardening: the daemon only needs to read its config and write the vault working tree.
# It also shells out to `hindsight-embed -p thoth memory recall/retain`, which connects to
# the separate thoth-hindsight.service over 127.0.0.1:9277 but may still touch client-side
# state under ~/.hindsight (and the embedded PG dir ~/.pg0). Under ProtectHome=read-only
# those writes would be denied, so both are punched through here as ReadWritePaths. They
# are pkm's own dirs; the daemon never spawns Hindsight (that is thoth-hindsight.service).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/pkm-vault /home/pkm/.thoth /home/pkm/.hindsight /home/pkm/.pg0

[Install]
WantedBy=multi-user.target
